# backend/docker/supervisor/restaurant-management.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
priority=10
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true

[program:queue-worker]
command=php /var/www/html/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/queue-worker.log
stderr_logfile=/var/www/html/storage/logs/queue-worker-error.log
